# pitally/bin/pitally.service

[Unit]
Description=Pitally server
#Wants=ntpd.service
#After=ntpd.service


[Service]
Type=simple
#WorkingDirectory=/opt/ethoscope-git/node_src/scripts
ExecStart=pitally.sh
RestartSec=5
Restart=always

[Install]
WantedBy=multi-user.target